# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

; Laravel Octane (Swoole worker — keeps app booted in memory between requests)
[program:octane]
command=php /var/www/html/artisan octane:start --server=swoole --host=127.0.0.1 --port=8000 --workers=8 --task-workers=4 --max-requests=1000
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/log/octane.log
stderr_logfile=/var/log/octane.err
stopsignal=TERM
stopwaitsecs=10

; Nginx — reverse proxy on :80 → Octane :8000 + static files directly
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true

; Queue worker — БЕЗ нього queued-джоби (листи замовлень, SMS/Viber,
; нотифікації ТТН) НІКОЛИ не обробляються (вічно висять у Redis).
[program:queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/log/queue.log
stderr_logfile=/var/log/queue.err
stopsignal=TERM
stopwaitsecs=15

; Scheduler — щохвилинний tick для розкладу (курси валют, обслуговування).
[program:scheduler]
command=sh -c "while true; do php /var/www/html/artisan schedule:run --no-interaction; sleep 60; done"
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/log/scheduler.log
stderr_logfile=/var/log/scheduler.err
